# templates/supervisord-services.ini
{% for service in global_supervisord_services %}
[program:{{ service.name }}]
command={{ service.command }}
environment={{ service.environment | default(supervisord_service_default_environment) }}
directory={{ service.directory | default(supervisord_service_default_directory) }}
user={{ service.user | default(supervisord_service_default_user) }}
autorestart={{ service.autorestart | default(supervisord_service_default_autorestart) }}
autostart={{ service.autostart | default(supervisord_service_default_autostart) }}
stopwaitsecs={{ service.stopwaitsecs | default(supervisord_service_default_stopwaitsecs) }}

{% endfor %}

{% for service in specific_supervisord_services %}
[program:{{ service.name }}]
command={{ service.command }}
environment={{ service.environment | default(supervisord_service_default_environment) }}
directory={{ service.directory | default(supervisord_service_default_directory) }}
user={{ service.user | default(supervisord_service_default_user) }}
autorestart={{ service.autorestart | default(supervisord_service_default_autorestart) }}
autostart={{ service.autostart | default(supervisord_service_default_autostart) }}
stopwaitsecs={{ service.stopwaitsecs | default(supervisord_service_default_stopwaitsecs) }}

{% endfor %}

